# webapp.service
[Unit]
Description=webapp
After=network.target
Wants=cloud-init.target

[Service]
WorkingDirectory=/opt/webappuser/webapp
ExecStart=/usr/bin/node /opt/webappuser/webapp/index.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=webapp
User=webappuser

[Install]
WantedBy=multi-user.target